# urender/files/etc/hotplug.d/iface/90-ip-collide
[ ifup = "$ACTION" ] && /usr/share/urender/ip-collide.uc
