# urender/files/etc/init.d/urender
#!/bin/sh /etc/rc.common

START=10

PROG=/usr/bin/ucode

boot() {
	[ -f /etc/urender/capabilities.json ] || /usr/share/urender/capabilities.uc
}
